# .crontab
#Update Reminder
0 * * * *	~/.config/sway/scripts/helpers/update_reminder.sh


# File generated by Crontablib the Tuesday, February 7, 2023 3:45:34 PM EST.
